# blereceive.service
[Unit]
Description=Run BLE scanner script
After=rfkill-unblockall.service bluetooth.target

[Service]
Type=forking
User=root
Group=root
PIDFile=/var/run/ibeacon-scan-upload.pid
ExecStart=/usr/bin/python2 /root/ble/c-receive/daemon_load.py start
TimeoutStartSec=10
#RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
