# blereceive.service
[Unit]
Description=Run BLE scanner script

[Service]
Type=forking
ExecStart=/root/ble/r-daemon-script.sh
TimeoutStartSec=10
#RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
